# squashfs-root/etc/uci-defaults/luci-i18n-swedish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.se=svenska
        commit luci
EOF

