# squashfs-root/etc/uci-defaults/luci-i18n-czech
#!/bin/sh
uci batch <<-EOF
	set luci.languages.cz=čeština
        commit luci
EOF

