# squashfs-root/etc/uci-defaults/luci-i18n-serbian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.rs=Српски
        commit luci
EOF

