# squashfs-root/etc/uci-defaults/luci-i18n-polish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.pl='Polskie'
	commit luci
EOF
