# squashfs-root/etc/init.d/smp
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=99

start() {
	(sleep 30 && /usr/sbin/smp.sh) &
}
