# squashfs-root/etc/hotplug.d/iface/99-sp_phy_down
#!/bin/ash

if [ "$ACTION" != "ifdown" ]; then
	exit 0
fi

if [ "$INTERFACE" == "wan" -o "$INTERFACE" == "lan" ]; then
	sp_log_info.sh -k net.phy.down -m "$INTERFACE:1"
fi

