# squashfs-root/usr/lib/opkg/info/base-files.list
/usr/sbin/update_whc_ob_bssid
/usr/sbin/killpidfile
/etc/config/miwifi
/sbin/setssid_fact
/etc/config/luci
/etc/hotplug.d/iface/70-vpn
/etc/openwrt_version
/etc/banner.failsafe
/usr/sbin/wanlinkprobe
/etc/init.d/powersaving
/lib/preinit/99_01_prepare_eeprom
/etc/resolv.conf
/usr/sbin/ban_admin
/usr/sbin/update_vas_conf.lua
/etc/init.d/done
/bin/config_generate
/usr/sbin/do_net_diagose.lua
/var
/etc/intercept_wlist
/usr/sbin/get_bdata_country
/etc/init.d/wifiradio
/lib/preinit/00_preinit.conf
/usr/sbin/autowancheck
/etc/init.d/system
/etc/preinit
/etc/init.d/network
/usr/sbin/iwpriv_macfilter.lua
/etc/hotplug.d/button/51-reset
/bin/reboot_safe.sh
/usr/sbin/guestwifi.sh
/etc/localtime
/usr/sbin/parentalctl.sh
/lib/upgrade/platform.sh
/etc/shells
/usr/sbin/miot_ip_conflict.lua
/bin/flash_check.sh
/lib64
/lib/preinit/70_initramfs_test
/etc/device_info
/usr/sbin/set_wps_state
/usr/lib64
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/usr/sbin/ip_conflict.sh
/usr/sbin/mivpn.sh
/etc/config/smartvpn
/lib/preinit/10_indicate_preinit
/etc/init.d/smp
/usr/libexec/login.sh
/usr/sbin/dhcp_apclient.sh
/usr/sbin/sysapi
/usr/sbin/supervisord
/etc/config/xqled
/usr/sbin/startscene.lua
/lib/preinit/05_set_preinit_iface
/usr/sbin/ddnsd
/etc/config/server_mapping
/usr/sbin/gensmartdns.sh
/etc/hotplug.d/button/02-mesh
/usr/sbin/logupload.lua
/sbin/chk_rebuild_miotcfg
/etc/config/macfilter
/usr/sbin/cr_ip_conflict.lua
/etc/config/system
/usr/sbin/smp.sh
/usr/sbin/sync_guest_bssid.lua
/sbin/wps
/etc/fstab
/etc/config/miwifi_host_mapping
/usr/lib/os-release
/etc/sysctl.conf
/usr/sbin/nslookupc
/lib/preinit/99_10_run_init
/usr/sbin/nslookupt
/etc/modules.d/92-mt7615e
/lib/preinit/98_10_mtk_failsafe_init
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/usr/sbin/crontab_record_camera.lua
/usr/sbin/kill_plugin_process.sh
/usr/lib/lua/json.lua
/usr/sbin/setchanauto.sh
/etc/shadow
/etc/rc.common
/sbin/setmac
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/usr/sbin/logcat
/etc/board.d/99-default_network
/usr/sbin/ap_mode.sh
/lib/preinit/39_mount_ubi_data
/lib/preinit/10_indicate_failsafe
/usr/share/xiaoqiang/xiaoqiang-defaults.txt
/etc/iwevent.d/01-feedpush.sh
/etc/hotplug.d/iface/10-phy_check.sh
/etc/os-release
/etc/board.d/02_network
/sbin/ethstatus
/usr/sbin/sysapi.firewall
/usr/sbin/jason.sh
/usr/sbin/re_mode.sh
/lib/upgrade/boardupgrade.sh
/usr/sbin/check_apcli_connected
/usr/sbin/iperf_script.lua
/etc/init.d/cgroup_init
/lib/dnsmasq/lease.d/10-dhcp-mitvbox.sh
/etc/config/diag
/etc/config/firewall
/lib/preinit/80_mount_root
/lib/preinit/81_urandom_seed
/sbin/chk_sta_re
/etc/hotplug.d/iface/60-wan
/etc/rc.local
/usr/sbin/set_wifi_auto_ap_mode.lua
/lib/functions/uci-defaults.sh
/usr/sbin/macfilterctl
/usr/sbin/bootinfo
/sbin/getmac
/usr/sbin/record_camera_by_ip.lua
/etc/ethers
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/preinit/90_mount_bind_etc
/sbin/hotplug-call
/usr/sbin/runt
/etc/crontabs/root
/lib/preinit/06_set_rps_sock_flow
/usr/sbin/getneighbor.sh
/usr/sbin/crontab_rom.sh
/etc/hotplug.d/iface/99-mtk-lro
/lib/preinit/03_preinit_do_mediatek.sh
/etc/config/country_mapping
/etc/services
/etc/config/vas_user
/lib/preinit/10_check_for_mtd
/lib/functions/system.sh
/bin/flash.sh
/usr/sbin/vasinfo_fw.sh
/lib/upgrade/libupgrade.sh
/bin/board_detect
/usr/sbin/dhcp_do_opt43_act.sh
/sbin/setmac_all
/etc/init.d/urandom_seed
/lib/preinit/39_mount_tmpfs
/etc/opkg/distfeeds.conf
/lib/preinit/31_restore_nvram
/etc/opkg/keys/bec8fa53cb562ca9
/etc/sysupgrade.conf
/usr/sbin/getCountryCode
/usr/sbin/proxy_thirdparty.sh
/lib/upgrade/keep.d/base-files
/usr/sbin/ecos_upgrade.lua
/etc/parentalctl/parentalctl.empty
/lib/preinit/40_run_failsafe_hook
/etc/group
/usr/sbin/startscene_crontab.lua
/etc/config/vas
/lib/preinit/30_failsafe_wait
/etc/uci-defaults/13_fix_group_user
/usr/sbin/showmacs
/bin/login.sh
/usr/sbin/config_post_ota.sh
/etc/hotplug.d/iface/97-clear-hnat
/etc/hotplug.d/net/00-sysctl
/lib/xq-misc/arch_function_lib.sh
/etc/hosts
/etc/wan_check/conned.d/40-ntpd
/lib/functions.sh
/usr/sbin/smartvpn.sh
/etc/remove.conf
/etc/custom_hosts
/etc/config/network
/usr/sbin/upcrash
/etc/config/region_mapping
/lib/xq-misc/phy_switch.sh
/etc/smartvpn/proxy.txt
/usr/sbin/vip_device_push_act.lua
/etc/config/parentalctl
/etc/init.d/timezone
/etc/init.d/autovpn
/usr/sbin/ntpsetclock
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/hwversion
/usr/sbin/secboot_upgrade_check.sh
/usr/share/xiaoqiang/xiaoqiang_version
/sbin/wifi
/lib/network/network_func_lib.sh
/lib/wifi/ralink_common.sh
/usr/sbin/led_ctl
/etc/profile
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/usr/share/xiaoqiang/xiaoqiang-reserved.txt
/etc/TZ
/etc/config/misc
/etc/hotplug.d/iface/98-mtk-vpn
/etc/init.d/umount
/lib/functions/service.sh
/etc/sysctl.d/10-default.conf
/usr/sbin/log_collection.sh
/etc/passwd
/usr/sbin/check_5g_state
/lib/preinit/02_default_set_state
/usr/sbin/speed_test.lua
/etc/mtab
/usr/sbin/dhcpc_do_opt43_act.sh
/etc/init.d/boot
/usr/sbin/phyhelper
/usr/sbin/config_pre_ota.sh
/usr/sbin/doStatPoints
/sbin/urandom_seed
/lib/functions/preinit.sh
/lib/mediatek.sh
/usr/sbin/syslog_upload.lua
/usr/sbin/checkupgrade.lua
/bin/ipcalc.sh
/lib/firewall/ipset_redirect.loader
/etc/diag.sh
/usr/sbin/killpid
/usr/sbin/gpio
/usr/sbin/nettb
