# squashfs-root/usr/lib/opkg/info/mtkhnat_util.list
/etc/init.d/mtkhnat
/sbin/mtkhnat
/etc/uci-defaults/99-firewall
/etc/config/mtkhnat
