# squashfs-root/usr/lib/opkg/info/xqled.list
/lib/xqled/xqled_api.sh
/usr/sbin/xqled
/etc/init.d/xqled
/lib/xqled/xqled_hal.sh
